# scripts/nudge-update.timer
[Unit]
Description=Daily Claude Code update for Nudge

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
