# scripts/nudge-update.service
[Unit]
Description=Update Claude Code for Nudge
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/nudge/scripts/update-claude-code.sh
